# containers/quadlet/warp-app.container
[Unit]
Description=WARP reservation system application
# warp-db.container runs PostgreSQL in the same pod; all containers share localhost.
# If you prefer an external PostgreSQL instead, remove these two lines and update
# WARP_DATABASE_ADDRESS below with the external host address.
Requires=warp-db.service
After=warp-db.service

[Container]
Pod=warp.pod
StartWithPod=true
ContainerName=warp-app

# Image is built by warp-app.build / warp-app-build.service.
Image=warp-app.build

# Database connection. If PostgreSQL runs in the same pod, use localhost.
# For an external host, set DATABASE_ADDRESS to its host[:port].
Environment=WARP_DATABASE_ADDRESS=localhost:5432
Environment=WARP_DATABASE_NAME=warp
Environment=WARP_DATABASE_USER=postgres
Secret=warp-db-password,type=env,target=WARP_DATABASE_PASSWORD

# Cookie signing key — injected from the podman secret "warp-secret-key".
Secret=warp-secret-key,type=env,target=WARP_SECRET_KEY

# UI language. Available: en, de, fr, es, pl.
Environment=WARP_LANGUAGE_FILE=i18n/en.js

# Uncomment and adjust any additional settings from CONFIGURATION.md:
# Environment=WARP_WEEKS_IN_ADVANCE=1
# Environment=WARP_OMITTED_WEEKDAYS=[5, 6]
# Environment=WARP_MIN_PASSWORD_LENGTH=8

# OpenID Connect (OIDC) authentication — see CONFIGURATION.md
# Environment=WARP_AUTH_OIDC=true
# Environment=WARP_OIDC_DISCOVERY_URL=https://idp.example.org/realms/warp/.well-known/openid-configuration
# Environment=WARP_OIDC_CLIENT_ID=warp
# Secret=warp-oidc-client-secret,type=env,target=WARP_OIDC_CLIENT_SECRET
# Environment=WARP_OIDC_GROUP_MAP=[ ["warp-allowed", null], [null, "Everyone"] ]
# Create the podman secret with: podman secret create warp-oidc-client-secret <file>

[Install]
WantedBy=multi-user.target default.target
