# containers/quadlet/warp-app.container
[Unit]
Description=WARP reservation system application
# warp-db.container runs PostgreSQL in the same pod; all containers share localhost.
# If you prefer an external PostgreSQL instead, remove these two lines and update
# WARP_DATABASE_ADDRESS below with the external host address.
Requires=warp-db.service
After=warp-db.service

[Container]
Pod=warp.pod
StartWithPod=true
ContainerName=warp-app

# Build the image locally from containers/Dockerfile and tag it warp:latest:
#   podman build -f containers/Dockerfile -t warp:latest .
# Once the image is published to the registry, switch this to the pullable
# reference (e.g. ghcr.io/<owner>/warp:latest) and add AutoUpdate=registry.
Image=warp:latest

Volume=warp-shared.volume:/run/warp

# Database connection. If PostgreSQL runs in the same pod, use localhost.
# For an external host, set DATABASE_ADDRESS to its host[:port].
Environment=WARP_DATABASE_ADDRESS=localhost:5432
Environment=WARP_DATABASE_NAME=warp
Environment=WARP_DATABASE_USER=postgres

Secret=warp-db-password
Environment=WARP_DATABASE_PASSWORD_FILE=/run/secrets/warp-db-password

# Cookie signing key — injected from the podman secret "warp-secret-key".
Secret=warp-secret-key
Environment=WARP_SECRET_KEY_FILE=/run/secrets/warp-secret-key

# UI language. Available: en, de, fr, es, pl.
Environment=WARP_LANGUAGE_FILE=i18n/en.json

[Install]
WantedBy=multi-user.target default.target
