# containers/quadlet/warp-app.container
[Unit]
Description=WARP reservation system application
# warp-db.container runs PostgreSQL in the same pod; all containers share localhost.
# If you prefer an external PostgreSQL instead, remove these two lines and update
# WARP_DATABASE_ADDRESS below with the external host address.
Requires=warp-db.service
After=warp-db.service

[Container]
Pod=warp.pod
StartWithPod=true
ContainerName=warp-app

# Application image, built and published by the GitHub Actions workflow
# (.github/workflows/containers.yml) to the GitHub Container Registry:
#   https://github.com/sebo-b/warp/pkgs/container/warp
# Pin a version tag (e.g. :v1.2.3) in production instead of :latest.
Image=ghcr.io/sebo-b/warp:latest
AutoUpdate=registry

Volume=warp-shared.volume:/run/warp

# Database connection. If PostgreSQL runs in the same pod, use localhost.
# For an external host, set DATABASE_ADDRESS to its host[:port].
Environment=WARP_DATABASE_ADDRESS=localhost:5432
Environment=WARP_DATABASE_NAME=warp
Environment=WARP_DATABASE_USER=postgres

Secret=warp-db-password
Environment=WARP_DATABASE_PASSWORD_FILE=/run/secrets/warp-db-password

# Cookie signing key — injected from the podman secret "warp-secret-key".
Secret=warp-secret-key
Environment=WARP_SECRET_KEY_FILE=/run/secrets/warp-secret-key

# UI language. Available: en, de, fr, es, pl. JSON array (the picker only
# renders when more than one is listed); fallback must be listed in LANGUAGES.
Environment=WARP_LANGUAGES="[\"en\",\"de\",\"fr\",\"es\",\"pl\"]"
Environment=WARP_DEFAULT_LANGUAGE=en

# uWSGI endpoints. The in-pod Caddy proxies over the HTTP unix socket on the
# shared /run/warp volume, so the binary uWSGI socket is disabled (empty value).
# The HTTP socket already defaults to /run/warp/uwsgi-http.sock in the app, so it
# is left commented below. For an external reverse proxy, set one of these to a
# host:port instead, e.g. WARPAPP_UWSGI_HTTP_SOCKET=0.0.0.0:8080 (see warp.pod
# option B).
#Environment="WARPAPP_UWSGI_HTTP_SOCKET=/run/warp/uwsgi-http.sock"
Environment="WARPAPP_UWSGI_SOCKET="

[Install]
WantedBy=multi-user.target default.target
